# root/etc/services.d/libresonic/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc jetty-runner "${LIBRE_HOME}"/libresonic.war
