# root/etc/services.d/libresonic/run
#!/usr/bin/with-contenv bash

jetty-runner "${CATALINA_HOME}"/libresonic.war
